# linknx-uci/files/etc/init.d/linknx-uci
#!/bin/sh /etc/rc.common
# Copyright (C) 2012 Patrick Grimm

START=99
USE_PROCD=1
BIN="/usr/sbin/linknx-uci.lua"

load_uci() {
	/usr/bin/linknxloaddata.lua
}

start_service() {
	load_uci
	procd_open_instance
	procd_set_param respawn
	procd_set_param stderr 1
	procd_set_param command "$BIN"
	procd_open_trigger
	procd_add_raw_trigger "interface.*" 2000 /etc/init.d/linknx-uci restart
	config_files="network linknx bacnet_ai bacnet_ao bacnet_av bacnet_bi bacnet_bo bacnet_bv bacnet_mi bacnet_mo bacnet_mv"
	for file in $config_files ; do
		if [ -f "/etc/config/$file" ] ; then
			procd_add_config_trigger "config.change" "$file" /etc/init.d/linknx-uci restart
		fi
	done
	procd_close_trigger
	procd_close_instance
}

reload_service() {
	load_uci
}

restart_service() {
	stop
	/etc/init.d/linknx restart
	start
}
