# linknx-uci/files/etc/init.d/linknx-uci
#!/bin/sh /etc/rc.common
# Copyright (C) 2012 Patrick Grimm

START=99

load_uci() {
	/usr/bin/linknxloaddata.lua
}

start() {
	load_uci
}

reload() {
	load_uci
}

restart() {
	/etc/init.d/linknx restart
	start
}

stop() {
	/etc/init.d/linknx stop
}
